# System/T/Alien/supervisord/supervisord.d/flask.ini
[program:flask_api-8004]
command=/data/Py2.7/bin/uwsgi /data/Py2.7/wwwroot/flask_api/config.ini
directory=/data/Py2.7/wwwroot/flask_api
autorestart=true
redirect_stderr=true
stdout_logfile=/data/logs/flask.log
loglevel=info
