# System/T/Alien/supervisord/supervisord.d/texas.ini
[program:Alien]
command         = /data/Py27/bin/python2.7 /data/Py27/wwwroot/Alien/main.py
process_name    = Alien
numprocs        = 1
numprocs_start  = 1
autostart       = true
autorestart     = true
redirect_stderr = true
stdout_logfile  = /data/logs/alien_stdout.log
stderr_logfile  = /data/logs/alien_stdout.log

[program:FlaskAPI]
command=/data/Py27/bin/uwsgi /data/Py27/wwwroot/flask_api/config.ini
directory=/data/Py27/wwwroot/flask_api
autorestart=true
redirect_stderr=true
stdout_logfile=/data/logs/flaskAPI_stdout.log
stderr_logfile=/data/logs/flaskAPI_stderr.log
loglevel=info
